# data/v12/profiles/profile_adversary_temp.veles.json
[
  {
    "id": "intrusion-set--9538b1a4-4120-4e2d-bf59-3b11fcab05a4",
    "name": "temp.veles",
    "aliases": "temp.veles|xenotime",
    "group_id": "G0088",
    "is_revoked": false,
    "profile": {
      "malware": {
        "count": 1,
        "items": [
          "triton"
        ]
      },
      "tools": {
        "count": 2,
        "items": [
          "mimikatz",
          "psexec"
        ]
      },
      "techniques": {
        "count": 17,
        "items": [
          "T1003.001",
          "T1021.001",
          "T1021.004",
          "T1027.005",
          "T1036.005",
          "T1053.005",
          "T1059.001",
          "T1070.004",
          "T1070.006",
          "T1074.001",
          "T1078",
          "T1133",
          "T1505.003",
          "T1546.012",
          "T1571",
          "T1583.003",
          "T1588.002"
        ]
      },
      "tactics": {
        "count": 10,
        "items": [
          "collection",
          "command-and-control",
          "credential-access",
          "defense-evasion",
          "execution",
          "initial-access",
          "lateral-movement",
          "persistence",
          "privilege-escalation",
          "resource-development"
        ]
      },
      "procedures": [
        "T1003.001::[TEMP.Veles] has used [Mimikatz] and a custom tool, SecHack, to harvest credentials.",
        "T1021.001::[TEMP.Veles] utilized RDP throughout an operation.",
        "T1021.004::[TEMP.Veles] has relied on encrypted SSH-based tunnels to transfer tools and for remote command/program execution.",
        "T1027.005::[TEMP.Veles] has modified files based on the open-source project cryptcat in an apparent attempt to decrease AV detection rates.",
        "T1036.005::[TEMP.Veles] has renamed files to look like legitimate files, such as Windows update files or Schneider Electric application files.",
        "T1053.005::[TEMP.Veles] has used scheduled task XML triggers.",
        "T1059.001::[TEMP.Veles] has used a publicly-available PowerShell-based tool, WMImplant.",
        "T1070.004::[TEMP.Veles] routinely deleted tools, logs, and other files after they were finished with them.",
        "T1070.006::[TEMP.Veles] used timestomping to modify the $STANDARD_INFORMATION attribute on tools.",
        "T1074.001::[TEMP.Veles] has created staging folders in directories that were infrequently used by legitimate users or processes.",
        "T1078::[TEMP.Veles] has used compromised VPN accounts.",
        "T1133::[TEMP.Veles] has used a VPN to persist in the victim environment.",
        "T1505.003::[TEMP.Veles] has planted Web shells on Outlook Exchange servers.",
        "T1546.012::[TEMP.Veles] has modified and added entries within HKEY_LOCAL_MACHINE\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\Image File Execution Options to maintain persistence.",
        "T1571::[TEMP.Veles] has used port-protocol mismatches on ports such as 443, 4444, 8531, and 50501 during C2.",
        "T1583.003::[TEMP.Veles] has used Virtual Private Server (VPS) infrastructure.",
        "T1588.002::[TEMP.Veles] has obtained and used tools such as [Mimikatz] and [PsExec]."
      ]
    },
    "description": "[TEMP.Veles](https://attack.mitre.org/groups/G0088) is a Russia-based threat group that has targeted critical infrastructure. The group has been observed utilizing [TRITON](https://attack.mitre.org/software/S0609), a malware framework designed to manipulate industrial safety systems.(Citation: FireEye TRITON 2019)(Citation: FireEye TEMP.Veles 2018)(Citation: FireEye TEMP.Veles JSON April 2019)"
  }
]
