# data/current/dataset/procedures/procedures_adversary_fox-kitten.json
[
  [
    "fox-kitten",
    "T1003.001",
    "LSASS Memory",
    "command-execution|os-api-execution|process-access|process-creation",
    "T1003.001::[Fox Kitten] has used prodump to dump credentials from LSASS."
  ],
  [
    "fox-kitten",
    "T1003.003",
    "NTDS",
    "command-execution|file-access",
    "T1003.003::[Fox Kitten] has used Volume Shadow Copy to access credential information from NTDS."
  ],
  [
    "fox-kitten",
    "T1021.001",
    "Remote Desktop Protocol",
    "logon-session-creation|network-connection-creation|network-traffic-flow|process-creation",
    "T1021.001::[Fox Kitten] has used RDP to log in and move laterally in the target environment."
  ],
  [
    "fox-kitten",
    "T1021.002",
    "SMB/Windows Admin Shares",
    "command-execution|logon-session-creation|network-connection-creation|network-share-access|network-traffic-flow",
    "T1021.002::[Fox Kitten] has used valid accounts to access SMB shares."
  ],
  [
    "fox-kitten",
    "T1021.004",
    "SSH",
    "logon-session-creation|network-connection-creation|process-creation",
    "T1021.004::[Fox Kitten] has used the PuTTY and Plink tools for lateral movement."
  ],
  [
    "fox-kitten",
    "T1021.005",
    "VNC",
    "logon-session-creation|network-connection-creation|process-creation",
    "T1021.005::[Fox Kitten] has installed TightVNC server and client on compromised servers and endpoints for lateral movement."
  ],
  [
    "fox-kitten",
    "T1027.010",
    "Command Obfuscation",
    "command-execution|file-metadata|script-execution",
    "T1027.010::[Fox Kitten] has base64 encoded scripts to avoid detection."
  ],
  [
    "fox-kitten",
    "T1036.004",
    "Masquerade Task or Service",
    "command-execution|scheduled-job-metadata|scheduled-job-modification|service-creation|service-metadata",
    "T1036.004::[Fox Kitten] has named the task for a reverse proxy lpupdate to appear legitimate."
  ],
  [
    "fox-kitten",
    "T1036.005",
    "Match Legitimate Name or Location",
    "file-metadata|image-metadata|process-metadata",
    "T1036.005::[Fox Kitten] has named binaries and configuration files svhost and dllhost respectively to appear legitimate."
  ],
  [
    "fox-kitten",
    "T1053.005",
    "Scheduled Task",
    "command-execution|file-modification|process-creation|scheduled-job-creation|windows-registry-key-creation",
    "T1053.005::[Fox Kitten] has used Scheduled Tasks for persistence and to load and execute a reverse proxy binary."
  ],
  [
    "fox-kitten",
    "T1059.001",
    "PowerShell",
    "command-execution|module-load|process-creation|process-metadata|script-execution",
    "T1059.001::[Fox Kitten] has used PowerShell scripts to access credential data."
  ],
  [
    "fox-kitten",
    "T1059.003",
    "Windows Command Shell",
    "command-execution|process-creation",
    "T1059.003::[Fox Kitten] has used cmd.exe likely as a password changing mechanism."
  ],
  [
    "fox-kitten",
    "T1087.001",
    "Local Account",
    "command-execution|file-access|group-enumeration|os-api-execution|process-creation",
    "T1087.001::[Fox Kitten] has accessed ntuser.dat and UserClass.dat on compromised hosts."
  ],
  [
    "fox-kitten",
    "T1087.002",
    "Domain Account",
    "command-execution|group-enumeration|network-traffic-content|os-api-execution|process-creation",
    "T1087.002::[Fox Kitten] has used the Softerra LDAP browser to browse documentation on service accounts."
  ],
  [
    "fox-kitten",
    "T1136.001",
    "Local Account",
    "command-execution|process-creation|user-account-creation",
    "T1136.001::[Fox Kitten] has created a local user account with administrator privileges."
  ],
  [
    "fox-kitten",
    "T1505.003",
    "Web Shell",
    "application-log-content|file-creation|file-modification|network-traffic-content|network-traffic-flow|process-creation",
    "T1505.003::[Fox Kitten] has installed web shells on compromised hosts to maintain access."
  ],
  [
    "fox-kitten",
    "T1546.008",
    "Accessibility Features",
    "command-execution|file-creation|file-modification|process-creation|windows-registry-key-modification",
    "T1546.008::[Fox Kitten] has used sticky keys to launch a command prompt."
  ],
  [
    "fox-kitten",
    "T1552.001",
    "Credentials In Files",
    "command-execution|file-access",
    "T1552.001::[Fox Kitten] has accessed files to gain valid credentials."
  ],
  [
    "fox-kitten",
    "T1555.005",
    "Password Managers",
    "command-execution|file-access|os-api-execution|process-access",
    "T1555.005::[Fox Kitten] has used scripts to access credential information from the KeePass database."
  ],
  [
    "fox-kitten",
    "T1560.001",
    "Archive via Utility",
    "command-execution|file-creation|process-creation",
    "T1560.001::[Fox Kitten] has used 7-Zip to archive data."
  ],
  [
    "fox-kitten",
    "T1585.001",
    "Social Media Accounts",
    "network-traffic-content|social-media",
    "T1585.001::[Fox Kitten] has used a Twitter account to communicate with ransomware victims."
  ]
]
