# data/current/dataset/procedures/procedures_adversary_menupass.json
[
  [
    "menupass",
    "T1003.002",
    "Security Account Manager",
    "command-execution|file-access|windows-registry-key-access",
    "T1003.002::[menuPass] has used a modified version of pentesting tools wmiexec.vbs and secretsdump.py to dump credentials."
  ],
  [
    "menupass",
    "T1003.003",
    "NTDS",
    "command-execution|file-access",
    "T1003.003::[menuPass] has used Ntdsutil to dump credentials."
  ],
  [
    "menupass",
    "T1003.004",
    "LSA Secrets",
    "command-execution|windows-registry-key-access",
    "T1003.004::[menuPass] has used a modified version of pentesting tools wmiexec.vbs and secretsdump.py to dump credentials."
  ],
  [
    "menupass",
    "T1021.001",
    "Remote Desktop Protocol",
    "logon-session-creation|network-connection-creation|network-traffic-flow|process-creation",
    "T1021.001::[menuPass] has used RDP connections to move across the victim network."
  ],
  [
    "menupass",
    "T1021.004",
    "SSH",
    "logon-session-creation|network-connection-creation|process-creation",
    "T1021.004::[menuPass] has used Putty Secure Copy Client (PSCP) to transfer data."
  ],
  [
    "menupass",
    "T1036.003",
    "Rename System Utilities",
    "command-execution|file-metadata|file-modification|process-metadata",
    "T1036.003::[menuPass] has renamed [certutil] and moved it to a different location on the system to avoid detection based on use of the tool."
  ],
  [
    "menupass",
    "T1036.005",
    "Match Legitimate Name or Location",
    "file-metadata|image-metadata|process-metadata",
    "T1036.005::[menuPass] has been seen changing malicious files to appear legitimate."
  ],
  [
    "menupass",
    "T1053.005",
    "Scheduled Task",
    "command-execution|file-modification|process-creation|scheduled-job-creation|windows-registry-key-creation",
    "T1053.005::[menuPass] has used a script (atexec.py) to execute a command on a target machine via Task Scheduler."
  ],
  [
    "menupass",
    "T1055.012",
    "Process Hollowing",
    "os-api-execution|process-access|process-modification",
    "T1055.012::[menuPass] has used process hollowing in iexplore.exe to load the [RedLeaves] implant."
  ],
  [
    "menupass",
    "T1056.001",
    "Keylogging",
    "driver-load|os-api-execution|windows-registry-key-modification",
    "T1056.001::[menuPass] has used key loggers to steal usernames and passwords."
  ],
  [
    "menupass",
    "T1059.001",
    "PowerShell",
    "command-execution|module-load|process-creation|process-metadata|script-execution",
    "T1059.001::[menuPass] uses [PowerSploit] to inject shellcode into PowerShell."
  ],
  [
    "menupass",
    "T1059.003",
    "Windows Command Shell",
    "command-execution|process-creation",
    "T1059.003::[menuPass] executes commands using a command-line interface and reverse shell. The group has used a modified version of pentesting script wmiexec.vbs to execute commands."
  ],
  [
    "menupass",
    "T1070.003",
    "Clear Command History",
    "command-execution|file-deletion|file-modification|user-account-authentication",
    "T1070.003::[menuPass] has used [Wevtutil] to remove PowerShell execution logs."
  ],
  [
    "menupass",
    "T1070.004",
    "File Deletion",
    "command-execution|file-deletion",
    "T1070.004::A [menuPass] macro deletes files after it has decoded and decompressed them."
  ],
  [
    "menupass",
    "T1074.001",
    "Local Data Staging",
    "command-execution|file-access|file-creation|windows-registry-key-modification",
    "T1074.001::[menuPass] stages data prior to exfiltration in multi-part archives, often saved in the Recycle Bin."
  ],
  [
    "menupass",
    "T1074.002",
    "Remote Data Staging",
    "command-execution|file-access|file-creation",
    "T1074.002::[menuPass] has staged data on remote MSP systems or other victim networks prior to exfiltration."
  ],
  [
    "menupass",
    "T1087.002",
    "Domain Account",
    "command-execution|group-enumeration|network-traffic-content|os-api-execution|process-creation",
    "T1087.002::[menuPass] has used the Microsoft administration tool csvde.exe to export Active Directory data."
  ],
  [
    "menupass",
    "T1090.002",
    "External Proxy",
    "network-connection-creation|network-traffic-content|network-traffic-flow",
    "T1090.002::[menuPass] has used a global service provider's IP as a proxy for C2 traffic from a victim."
  ],
  [
    "menupass",
    "T1204.002",
    "Malicious File",
    "file-creation|process-creation",
    "T1204.002::[menuPass] has attempted to get victims to open malicious files such as Windows Shortcuts (.lnk) and/or Microsoft Office documents, sent via email as part of spearphishing campaigns."
  ],
  [
    "menupass",
    "T1218.004",
    "InstallUtil",
    "command-execution|process-creation",
    "T1218.004::[menuPass] has used InstallUtil.exe to execute malicious software."
  ],
  [
    "menupass",
    "T1553.002",
    "Code Signing",
    "file-metadata",
    "T1553.002::[menuPass] has resized and added data to the certificate table to enable the signing of modified files with legitimate signatures."
  ],
  [
    "menupass",
    "T1560.001",
    "Archive via Utility",
    "command-execution|file-creation|process-creation",
    "T1560.001::[menuPass] has compressed files before exfiltration using TAR and RAR."
  ],
  [
    "menupass",
    "T1566.001",
    "Spearphishing Attachment",
    "application-log-content|file-creation|network-traffic-content|network-traffic-flow",
    "T1566.001::[menuPass] has sent malicious Office documents via email as part of spearphishing campaigns as well as executables disguised as documents."
  ],
  [
    "menupass",
    "T1568.001",
    "Fast Flux DNS",
    "network-connection-creation|network-traffic-flow",
    "T1568.001::[menuPass] has used dynamic DNS service providers to host malicious domains."
  ],
  [
    "menupass",
    "T1574.001",
    "DLL Search Order Hijacking",
    "file-creation|file-modification|module-load",
    "T1574.001::[menuPass] has used DLL search order hijacking."
  ],
  [
    "menupass",
    "T1574.002",
    "DLL Side-Loading",
    "file-creation|file-modification|module-load|process-creation",
    "T1574.002::[menuPass] has used DLL side-loading to launch versions of Mimikatz and PwDump6 as well as [UPPERCUT]."
  ],
  [
    "menupass",
    "T1583.001",
    "Domains",
    "active-dns|domain-registration|passive-dns",
    "T1583.001::[menuPass] has registered malicious domains for use in intrusion campaigns."
  ],
  [
    "menupass",
    "T1588.002",
    "Tool",
    "malware-metadata",
    "T1588.002::[menuPass] has used and modified open-source tools like [Impacket], [Mimikatz], and [pwdump]."
  ]
]
