# data/v12/procedures/procedures_adversary_apt32.json
[
  [
    "apt32",
    "T1003.001",
    "LSASS Memory",
    "command-execution|os-api-execution|process-access|process-creation",
    "T1003.001::[APT32] used Mimikatz and customized versions of Windows Credential Dumper to harvest credentials."
  ],
  [
    "apt32",
    "T1021.002",
    "SMB/Windows Admin Shares",
    "command-execution|logon-session-creation|network-connection-creation|network-share-access|network-traffic-flow",
    "T1021.002::[APT32] used [Net] to use Windows' hidden network shares to copy their tools to remote machines for execution."
  ],
  [
    "apt32",
    "T1027.001",
    "Binary Padding",
    "file-metadata",
    "T1027.001::[APT32] includes garbage code to mislead anti-malware software and researchers."
  ],
  [
    "apt32",
    "T1036.003",
    "Rename System Utilities",
    "command-execution|file-metadata|file-modification|process-metadata",
    "T1036.003::[APT32] has moved and renamed pubprn.vbs to a .txt file to avoid detection."
  ],
  [
    "apt32",
    "T1036.004",
    "Masquerade Task or Service",
    "command-execution|scheduled-job-metadata|scheduled-job-modification|service-creation|service-metadata",
    "T1036.004::[APT32] has used hidden or non-printing characters to help masquerade service names, such as appending a Unicode no-break space character to a legitimate service name. [APT32] has also impersonated the legitimate Flash installer file name \"install_flashplayer.exe\"."
  ],
  [
    "apt32",
    "T1036.005",
    "Match Legitimate Name or Location",
    "file-metadata|image-metadata|process-metadata",
    "T1036.005::[APT32] has renamed a NetCat binary to kb-10233.exe to masquerade as a Windows update. [APT32] has also renamed a Cobalt Strike beacon payload to install_flashplayers.exe."
  ],
  [
    "apt32",
    "T1048.003",
    "Exfiltration Over Unencrypted Non-C2 Protocol",
    "command-execution|file-access|network-connection-creation|network-traffic-content|network-traffic-flow",
    "T1048.003::[APT32]'s backdoor can exfiltrate data by encoding it in the subdomain field of DNS packets."
  ],
  [
    "apt32",
    "T1053.005",
    "Scheduled Task",
    "command-execution|file-modification|process-creation|scheduled-job-creation|windows-registry-key-creation",
    "T1053.005::[APT32] has used scheduled tasks to persist on victim systems."
  ],
  [
    "apt32",
    "T1056.001",
    "Keylogging",
    "driver-load|os-api-execution|windows-registry-key-modification",
    "T1056.001::[APT32] has abused the PasswordChangeNotify to monitor for and capture account password changes."
  ],
  [
    "apt32",
    "T1059.001",
    "PowerShell",
    "command-execution|module-load|process-creation|process-metadata|script-execution",
    "T1059.001::[APT32] has used PowerShell-based tools, PowerShell one-liners, and shellcode loaders for execution."
  ],
  [
    "apt32",
    "T1059.003",
    "Windows Command Shell",
    "command-execution|process-creation",
    "T1059.003::[APT32] has used cmd.exe for execution."
  ],
  [
    "apt32",
    "T1059.005",
    "Visual Basic",
    "command-execution|module-load|process-creation|script-execution",
    "T1059.005::[APT32] has used macros, COM scriptlets, and VBS scripts."
  ],
  [
    "apt32",
    "T1059.007",
    "JavaScript",
    "command-execution|module-load|process-creation|script-execution",
    "T1059.007::[APT32] has used JavaScript for drive-by downloads and C2 communications."
  ],
  [
    "apt32",
    "T1070.001",
    "Clear Windows Event Logs",
    "command-execution|file-deletion|os-api-execution",
    "T1070.001::[APT32] has cleared select event log entries."
  ],
  [
    "apt32",
    "T1070.004",
    "File Deletion",
    "command-execution|file-deletion",
    "T1070.004::[APT32]'s macOS backdoor can receive a “delete” command."
  ],
  [
    "apt32",
    "T1070.006",
    "Timestomp",
    "file-metadata|file-modification",
    "T1070.006::[APT32] has used scheduled task raw XML with a backdated timestamp of June 2, 2016. The group has also set the creation time of the files dropped by the second stage of the exploit to match the creation time of kernel32.dll. Additionally, [APT32] has used a random value to modify the timestamp of the file storing the clientID."
  ],
  [
    "apt32",
    "T1071.001",
    "Web Protocols",
    "network-traffic-content|network-traffic-flow",
    "T1071.001::[APT32] has used JavaScript that communicates over HTTP or HTTPS to attacker controlled domains to download additional frameworks. The group has also used downloaded encrypted payloads over HTTP."
  ],
  [
    "apt32",
    "T1071.003",
    "Mail Protocols",
    "network-traffic-content|network-traffic-flow",
    "T1071.003::[APT32] has used email for C2 via an Office macro."
  ],
  [
    "apt32",
    "T1078.003",
    "Local Accounts",
    "logon-session-creation|logon-session-metadata|user-account-authentication",
    "T1078.003::[APT32] has used legitimate local admin account credentials."
  ],
  [
    "apt32",
    "T1087.001",
    "Local Account",
    "command-execution|file-access|group-enumeration|os-api-execution|process-creation",
    "T1087.001::[APT32] enumerated administrative users using the commands net localgroup administrators."
  ],
  [
    "apt32",
    "T1204.001",
    "Malicious Link",
    "file-creation|network-connection-creation|network-traffic-content",
    "T1204.001::[APT32] has lured targets to download a Cobalt Strike beacon by including a malicious link within spearphishing emails."
  ],
  [
    "apt32",
    "T1204.002",
    "Malicious File",
    "file-creation|process-creation",
    "T1204.002::[APT32] has attempted to lure users to execute a malicious dropper delivered via a spearphishing attachment."
  ],
  [
    "apt32",
    "T1216.001",
    "PubPrn",
    "command-execution|process-creation|script-execution",
    "T1216.001::[APT32] has used PubPrn.vbs within execution scripts to execute malware, possibly bypassing defenses."
  ],
  [
    "apt32",
    "T1218.005",
    "Mshta",
    "command-execution|file-creation|network-connection-creation|process-creation",
    "T1218.005::[APT32] has used mshta.exe for code execution."
  ],
  [
    "apt32",
    "T1218.010",
    "Regsvr32",
    "command-execution|module-load|network-connection-creation|process-creation",
    "T1218.010::[APT32] created a [Scheduled Task/Job] that used regsvr32.exe to execute a COM scriptlet that dynamically downloaded a backdoor and injected it into memory. The group has also used regsvr32 to run their backdoor."
  ],
  [
    "apt32",
    "T1218.011",
    "Rundll32",
    "command-execution|file-metadata|module-load|process-creation",
    "T1218.011::[APT32] malware has used rundll32.exe to execute an initial infection process."
  ],
  [
    "apt32",
    "T1222.002",
    "Linux and Mac File and Directory Permissions Modification",
    "command-execution|file-metadata|process-creation",
    "T1222.002::[APT32]'s macOS backdoor changes the permission of the file it wants to execute to 755."
  ],
  [
    "apt32",
    "T1505.003",
    "Web Shell",
    "application-log-content|file-creation|file-modification|network-traffic-content|network-traffic-flow|process-creation",
    "T1505.003::[APT32] has used Web shells to maintain access to victim websites."
  ],
  [
    "apt32",
    "T1543.003",
    "Windows Service",
    "command-execution|driver-load|os-api-execution|process-creation|service-creation|service-modification|windows-registry-key-creation|windows-registry-key-modification",
    "T1543.003::[APT32] modified Windows Services to ensure PowerShell scripts were loaded on the system. [APT32] also creates a Windows service to establish persistence."
  ],
  [
    "apt32",
    "T1547.001",
    "Registry Run Keys / Startup Folder",
    "command-execution|file-modification|process-creation|windows-registry-key-creation|windows-registry-key-modification",
    "T1547.001::[APT32] established persistence using Registry Run keys, both to execute PowerShell and VBS scripts as well as to execute their backdoor directly."
  ],
  [
    "apt32",
    "T1550.002",
    "Pass the Hash",
    "active-directory-credential-request|logon-session-creation|user-account-authentication",
    "T1550.002::[APT32] has used pass the hash for lateral movement."
  ],
  [
    "apt32",
    "T1550.003",
    "Pass the Ticket",
    "active-directory-credential-request|logon-session-creation|user-account-authentication",
    "T1550.003::[APT32] successfully gained remote access by using pass the ticket."
  ],
  [
    "apt32",
    "T1552.002",
    "Credentials in Registry",
    "command-execution|process-creation|windows-registry-key-access",
    "T1552.002::[APT32] used Outlook Credential Dumper to harvest credentials stored in Windows registry."
  ],
  [
    "apt32",
    "T1564.001",
    "Hidden Files and Directories",
    "command-execution|file-creation|file-metadata|process-creation",
    "T1564.001::[APT32]'s macOS backdoor hides the clientID file via a chflags function."
  ],
  [
    "apt32",
    "T1564.003",
    "Hidden Window",
    "command-execution|file-modification|process-creation|script-execution",
    "T1564.003::[APT32] has used the WindowStyle parameter to conceal [PowerShell] windows."
  ],
  [
    "apt32",
    "T1564.004",
    "NTFS File Attributes",
    "command-execution|file-metadata|file-modification|os-api-execution",
    "T1564.004::[APT32] used NTFS alternate data streams to hide their payloads."
  ],
  [
    "apt32",
    "T1566.001",
    "Spearphishing Attachment",
    "application-log-content|file-creation|network-traffic-content|network-traffic-flow",
    "T1566.001::[APT32] has sent spearphishing emails with a malicious executable disguised as a document or spreadsheet."
  ],
  [
    "apt32",
    "T1566.002",
    "Spearphishing Link",
    "application-log-content|network-traffic-content|network-traffic-flow",
    "T1566.002::[APT32] has sent spearphishing emails containing malicious links."
  ],
  [
    "apt32",
    "T1569.002",
    "Service Execution",
    "command-execution|process-creation|service-creation|windows-registry-key-modification",
    "T1569.002::[APT32]'s backdoor has used Windows services as a way to execute its malicious payload."
  ],
  [
    "apt32",
    "T1574.002",
    "DLL Side-Loading",
    "file-creation|file-modification|module-load|process-creation",
    "T1574.002::[APT32] ran legitimately-signed executables from Symantec and McAfee which load a malicious DLL. The group also side-loads its backdoor by dropping a library and a legitimate, signed executable (AcroTranscoder)."
  ],
  [
    "apt32",
    "T1583.001",
    "Domains",
    "active-dns|domain-registration|passive-dns",
    "T1583.001::[APT32] has set up and operated websites to gather information and deliver malware."
  ],
  [
    "apt32",
    "T1583.006",
    "Web Services",
    "response-content",
    "T1583.006::[APT32] has set up Dropbox, Amazon S3, and Google Drive to host malicious downloads."
  ],
  [
    "apt32",
    "T1585.001",
    "Social Media Accounts",
    "network-traffic-content|social-media",
    "T1585.001::[APT32] has set up Facebook pages in tandem with fake websites."
  ],
  [
    "apt32",
    "T1588.002",
    "Tool",
    "malware-metadata",
    "T1588.002::[APT32] has obtained and used tools such as [Mimikatz] and [Cobalt Strike], and a variety of other open-source tools from GitHub."
  ],
  [
    "apt32",
    "T1589.002",
    "Email Addresses",
    "network-traffic-content",
    "T1589.002::[APT32] has collected e-mail addresses for activists and bloggers in order to target them with spyware."
  ],
  [
    "apt32",
    "T1598.003",
    "Spearphishing Link",
    "application-log-content|network-traffic-content|network-traffic-flow",
    "T1598.003::[APT32] has used malicious links to direct users to web pages designed to harvest credentials."
  ],
  [
    "apt32",
    "T1608.001",
    "Upload Malware",
    "response-content",
    "T1608.001::[APT32] has hosted malicious payloads in Dropbox, Amazon S3, and Google Drive for use during targeting."
  ],
  [
    "apt32",
    "T1608.004",
    "Drive-by Target",
    "response-content",
    "T1608.004::[APT32] has stood up websites containing numerous articles and content scraped from the Internet to make them appear legitimate, but some of these pages include malicious JavaScript to profile the potential victim or infect them via a fake software update."
  ]
]
